# bot/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:process1]
command=python ./main.py
autostart=true

[program:process2]
command=python ./scheduled_tasks.py
autostart=true


